# deployment/systemd/cp2b-backend.service
[Unit]
Description=CP2B Backend API Server
Documentation=https://github.com/aikiesan/cp2bfun
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/cp2b/backend
EnvironmentFile=/var/www/cp2b/backend/.env

# Start the backend server
ExecStart=/usr/bin/node src/index.js

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitInterval=200
StartLimitBurst=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/cp2b/backend/uploads

# Environment
Environment=NODE_ENV=production

# Standard output and error logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cp2b-backend

# Resource limits
LimitNOFILE=65535
TasksMax=4096

[Install]
WantedBy=multi-user.target
